# assistant_backend/scripts/com.maldaze.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maldaze.backend</string>

    <key>ProgramArguments</key>
    <array>
        <string>UVICORN_PATH_PLACEHOLDER</string>
        <string>src.main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8765</string>
    </array>

    <key>WorkingDirectory</key>
    <string>BACKEND_DIR_PLACEHOLDER</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>GEMINI_API_KEY</key>
        <string>GEMINI_API_KEY_PLACEHOLDER</string>
        <key>DB_PATH</key>
        <string>DB_PATH_PLACEHOLDER</string>
        <key>PLAN_MD_PATH</key>
        <string>PLAN_MD_PATH_PLACEHOLDER</string>
        <key>PORT</key>
        <string>8765</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>LOG_DIR_PLACEHOLDER/maldaze-backend.log</string>

    <key>StandardErrorPath</key>
    <string>LOG_DIR_PLACEHOLDER/maldaze-backend-error.log</string>
</dict>
</plist>
